# crontab
0 12 * * * ruby /path/to/dosign.rb >> /path/to/log.txt
